# profile.d/fun.symlink
# hehe
alias ny='telnet miku.acm.uiuc.edu'
alias rr='curl -L http://bit.ly/10hA8iC | bash'
alias sw='telnet towel.blinkenlights.nl'

if [[ $(uname -s) == "Darwin" ]]; then
    source ~/.profile.osx.d/fun
elif [[ $(uname -s) == "Linux" ]]; then
    source ~/.profile.linux.d/fun
fi
